# scripts/crontab/3d.sh
#!/bin/sh

cd /home/wwwroot/yaf/
/usr/local/php7/bin/php cli.php request_uri="/weapp/jobnew/fc3d" >> ./logs/job/3d.log