# scripts/crontab/plw.sh
#!/bin/sh

cd /home/wwwroot/yaf/
/usr/local/php7/bin/php cli.php request_uri="/jobnew/plw" >> ./logs/Job/plw.log