# scripts/crontab/qxc.sh
#!/bin/sh

cd /home/wwwroot/yaf/
/usr/local/php7/bin/php cli.php request_uri="/jobnew/qxc" >> ./application/logs/Job/qxc.log